# cron.d/geokrety

MAILTO="contact+cron@geokretymap.org"

# Admin GeoKrety Databases (DEV)
*/15 *   * * * gkmap curl -s http://api.geokretymap.org/update2.php > /dev/null && /home/gkmap/web/admin-scripts/rrd/graph_rrd.sh > /dev/null
36   *   * * * gkmap curl -s http://api.geokretymap.org/optimize.php > /dev/null
20   */3 * * * gkmap curl -s http://api.geokretymap.org/update-details.php > /dev/null
35   */3 * * * gkmap curl -s http://api.geokretymap.org/backup.php > /dev/null
