# cron.d/geokrety-dev

MAILTO="contact+cron@geokretymap.org"

# Admin GeoKrety Databases (PROD)
*/15 *   * * * gkmap curl -s http://api.geokretymap.org/update2.php > /dev/null
36   *   * * * gkmap curl -s http://api.geokretymap.org/optimize.php > /dev/null
20   */3 * * * gkmap curl -s http://api.geokretymap.org/update-details.php > /dev/null
35   */3 * * * gkmap curl -s http://api.geokretymap.org/backup.php > /dev/null
