# cron.d/geokrety

MAILTO="contact+cron@geokretymap.org"

# Admin GeoKrety Databases (DEV)
30   2   * * * gkmap-dev curl -s http://api.dev.geokretymap.org/update2.php > /dev/null
35   2   * * * gkmap-dev curl -s http://api.dev.geokretymap.org/update-details.php > /dev/null
40   2   * * * gkmap-dev curl -s http://api.dev.geokretymap.org/optimize.php > /dev/null
45   2   * * * gkmap-dev curl -s http://api.dev.geokretymap.org/backup.php > /dev/null
